# packages/xcompmgr/init/xcompmgr.conf
description  "configure and launch xcompmgr"
author       "GeeXboX"

start on started xorg and runlevel geexbox
stop on stopping xorg
respawn

exec xcompmgr -n
