# packages/misc/nzbget/unit/nzbget.service
[Service]
# Even if you set up the Daemon User in nzbget.conf, you will still
# need to set that name here. Otherwise, nzbget will be run as root
# User=nzbget

# Use the 'LockFile' specified in nzbget.conf
PIDFile=/root/nzbget/nzbget.lock

Type=forking
ExecStart=/usr/bin/nzbget -D
ExecStop=/usr/bin/nzbget -Q
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process

[Install]
WantedBy=multi-user.target
