# packages/busybox/init/acpid.conf
description	"start the ACPI daemon subsystem"
author		"Davide Cavalca <davide@geexbox.org>"

start on started klogd

exec acpid
