# packages/busybox/init/acpid.conf
description  "start the ACPI daemon subsystem"
author       "GeeXboX"

start on started klogd

exec acpid -d
