# packages/busybox/init/crond.conf
description  "start cron daemon"
author       "GeeXboX"

start on started klogd
respawn

pre-start exec mkdir -p /var/spool/cron/crontabs

exec crond -fS
