# packages/busybox/init/poweroff.conf
description	"shutdown the system"
author		"GeeXboX"

start on shutdown REBOOT=no
task

exec poweroff
