# packages/busybox/init/reboot.conf
description	"reboot the system"
author		"GeeXboX"

start on shutdown REBOOT=yes
task

exec reboot
