# packages/busybox/init/shell2.conf
description  "Start debug shell"
author       "OpenBricks"

start on startup
respawn

env ENV

exec /sbin/getty -n -l /bin/sh 38400 tty2
