# packages/dropbear/unit/sshd@.service
[Unit]
Description=SSH Per-Connection Server
After=syslog.target

[Service]
ExecStartPre=/usr/bin/dropbearkeygen
ExecStart=/usr/bin/dropbear -i
StandardInput=socket
