# packages/dropbear/unit/sshd.socket
[Unit]
Description=SSH socket

[Socket]
ListenStream=22
Accept=yes

[Install]
WantedBy=sockets.target
