# config/platforms/arm/bcm2708/packages/xbmc/unit/silent.service
[Unit]
Description=Silent messages
After=xbmc.service

[Service]
Type=oneshot
ExecStart=/usr/bin/silent.sh

[Install]
WantedBy=xbmc.service
