# config/platforms/arm/bcm2708/packages/xbmc/unit/xbmc-setup.service
[Unit]
Description=XBMC Early Setup
Before=xbmc.service
After=graphical.target

[Service]
Type=oneshot
# Environment=HOME=/root
ExecStart=/usr/lib/xbmc/xbmc-default-settings.sh

[Install]
WantedBy=xbmc.service
