# packages/libs/rtorrent/unit/rtorrent-on-usb-disk@.service
[Unit]
Description=Search /rtorrent on usb disk
BindsTo=dev-%i.device

[Service]
ExecStart=/usr/bin/check-rtorrent-usb-disk.sh %I
ExecStop=/usr/bin/check-rtorrent-usb-disk-stop.sh %I
RemainAfterExit=yes
