# packages/enna/unit/enna.service
[Unit]
Description=Enna Media Center

[Service]
ExecStart=/usr/bin/runenna
StandardOutput=syslog
Restart=on-failure

[Install]
WantedBy=enna.target
