# packages/misc/irmplircd/udev/80-irmp.rules
#KERNEL=="hidraw0", ATTRS{idVendor}=="16c0", ATTRS{idProduct}=="05df", RUN+="/bin/mkdir /var/run/lirc", RUN+="/usr/local/bin/irmplircd -t /home/mark/irmplircd/remoto.map /dev/%k", RUN+="/bin/ln -s /var/run/lirc/lircd /dev/lircd"
KERNEL=="hidraw?", ATTRS{idVendor}=="16c0", ATTRS{idProduct}=="05df", RUN+="/bin/mkdir /var/run/lirc", RUN+="/usr/bin/irmplircd -t /etc/vdr/remoto.map /dev/%k", RUN+="/bin/ln -s /var/run/lirc/lircd /dev/lircd"
