# packages/eventlircd/unit/eventlircd.service
[Unit]
Description=Eventlircd daemon
After=udev.service
Before=xorg.service

[Service]
EnvironmentFile=/etc/eventlircd.conf
Type=forking
ExecStart=/usr/sbin/eventlircd $EV_OPTS

[Install]
WantedBy=multi-user.target

