# packages/net/udhcpc/unit/udhcpc.service
[Unit]
Description=Connection service (udhcpc)
After=syslog.target network-link.service
ConditionPathExists=/run/geexbox/network/udhcpc

[Service]
ExecStart=/usr/lib/udhcpc/net-udhcpc
RemainAfterExit=yes

[Install]
WantedBy=network.target
