# packages/net/bluez/config/bluetooth.rules
# Run helper every time a Bluetooth device appears
# On remove actions, bluetoothd should go away by itself
ACTION=="add", KERNEL=="hci[0-9]", RUN+="/usr/bin/hciconfig %k up"
