# packages/misc/netceiver/unit/mcli.service
[Unit]
Description=dvbloop&mcli

[Service]
ExecStart=/usr/sbin/dvbloop
Restart=always

[Install]
WantedBy=multi-user.target
