# packages/mex/unit/mex.service
[Unit]
Description=Media Explorer Media Center

[Service]
Environment=DISPLAY=:0.0 HOME=/root
ExecStart=/usr/bin/mex -f
Restart=always

[Install]
WantedBy=mex.target
