# packages/MPlayer/init/mplayer-config.conf
description  "configure MPlayer"
author       "OpenBricks"

start on started mplayer-fonts and started mplayer-audio and started mplayer-video
