# packages/bftpd/unit/bftpd-config.service
[Unit]
Description=FTP Server Configuration
ConditionPathExists=!/etc/bftpdpwd
Before=bftpd.service

[Service]
Type=oneshot
EnvironmentFile=/etc/ftp
ExecStart=/bin/sh -c 'echo "$USERNAME $PASSWORD root /" > /etc/bftpdpwd'
RemainAfterExit=yes

[Install]
WantedBy=bftpd.service
