# packages/multimedia/libomxil-bellagio/unit/omx-bellagio.service
[Unit]
Description=Register OpenMAX Bellagio plugins

[Service]
Type=oneshot
ExecStart=/usr/bin/omxregister-bellagio
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
