# packages/net/automounter-nfs/unit/automounter-nfs.service
[Unit]
Description=Mount NFS resources
ConditionPathExists=/etc/nfs
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/automount_nfs.sh
RemainAfterExit=yes

[Install]
#WantedBy=multi-user.target
WantedBy=network.target
