# packages/misc/haserl/unit/httpd.service
[Unit]
Description=Httpd Service

[Service]
ExecStart=/usr/sbin/httpd -f -v
Restart=always

[Install]
WantedBy=multi-user.target
