# packages/xbmc-pvr/unit/silent.service
[Unit]
Description=Silent messages
After=xbmc-pvr.service

[Service]
Type=oneshot
ExecStart=/usr/bin/silent.sh

[Install]
WantedBy=xbmc-pvr.service
