# packages/xbmc-pvr/unit/xbmc-pvr.service
[Unit]
Description=XBMC Media Center

[Service]
Environment=DISPLAY=:0.0 SDL_MOUSE_RELATIVE=0 HOME=/root
ExecStart=/usr/lib/xbmc/xbmc.bin --standalone -fs -l /var/run/lirc/lircd
Restart=always

[Install]
WantedBy=xbmc-pvr.target
