# packages/xbmc-pvr/unit/xbmc.service
[Unit]
Description=XBMC Media Center

[Service]
Environment=DISPLAY=:0.0
ExecStart=/usr/lib/xbmc/xbmc.bin
Restart=always

[Install]
WantedBy=xbmc.target
