# packages/vdr/unit/vdr.service
[Unit]
Description=Video Disk Recorder
After=lircd.service pvr.service
ConditionPathExists=/run/geexbox/pvr/vdr

[Service]
ExecStart=/usr/bin/runvdr
Restart=always

[Install]
WantedBy=multi-user.target
