# packages/vdr/unit/vdr.service
[Unit]
Description=Video Disk Recorder
After=lircd.service

[Service]
ExecStart=/usr/bin/runvdr
Restart=always

[Install]
WantedBy=multi-user.target
