# packages/vdr/unit/vdr.service
[Unit]
Description=Video Disk Recorder
After=lircd.service
ConditionPathExists=/dev/dvb/adapter0
ConditionPathExists=/etc/vdr/channels.conf

[Service]
ExecStart=/usr/bin/runvdr
Restart=always

[Install]
WantedBy=multi-user.target
