# packages/vdr/unit/vdr.service
[Unit]
Description=Video Disk Recorder
After=lircd.service

[Service]
EnvironmentFile=/etc/pvr.conf
ExecStartPre=/usr/bin/test ${BACKEND} = vdr
ExecStart=/usr/bin/runvdr
Restart=always

[Install]
WantedBy=multi-user.target
