# packages/LVM/init/lvm.conf
description  "mount LVM logical volumes"
author       "GeeXboX"

start on startup
task

exec automount_lvm
