# packages/misc/couchpotato/unit/couchpotato.service
[Unit]
Description=CouchPotato application instance
After=network.target

[Service]
ExecStart=/bin/python /srv/couchpotato/CouchPotato.py --daemon 
GuessMainPID=no
Type=forking
User=root
Group=root

[Install]
WantedBy=multi-user.target
