# packages/udev/config/rules.d/10_firmware.rules
ACTION=="add",  SUBSYSTEM=="firmware",  FIRMWARE="*",   RUN+="/sbin/firmware_helper"
