# packages/udev/config/rules.d/20-custom.rules
# custom rules

# cdrom
SUBSYSTEM=="block", KERNEL=="sr[0-9]*", SYMLINK+="scd%n", SYMLINK+="cdrom dvd cdroms/%k"
SUBSYSTEM=="block", KERNEL=="hd*", SUBSYSTEMS=="ide", ATTRS{media}=="cdrom", SYMLINK+="cdrom dvd cdroms/%k"

# lirc
KERNEL=="lirc[0-9]*", SYMLINK="lirc"

# video4linux
SUBSYSTEM=="video4linux", RUN+="/usr/bin/killall -TERM vdr"
