# packages/udev/config/rules.d/30_dvb.rules
# Video devices, group dvb devices under /dev/dvb
SUBSYSTEM!="dvb", GOTO="dvb_end"
IMPORT{program}="dvb_device_name --export %k"
ENV{DVB_ADAPTER}=="?*", ENV{DVB_DEV}=="?*", \
        NAME="dvb/adapter$env{DVB_ADAPTER}/$env{DVB_NAME}", \
        RUN+="/bin/sh -c 'echo > /var/dvbcard'"
LABEL="dvb_end"
