# packages/udev/config/rules.d/90_modprobe.rules
# Manage drivers that match kernel-supplied alias
ACTION=="add",  ENV{MODALIAS}=="?*",   RUN+="/sbin/modprobe $env{MODALIAS}"
ACTION=="remove",  ENV{MODALIAS}=="?*",   RUN+="/sbin/modprobe -r $env{MODALIAS}"
