# packages/xorg-server/unit/video-card.service
[Unit]
Description=Video card detection
After=lspci.service

[Service]
Type=oneshot
EnvironmentFile=/etc/video
ExecStart=/usr/lib/xorg/video-card-detect
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
