# packages/xorg-server/unit/xorg.service
[Unit]
Description=Xorg server
After=video-card.service input.service

[Service]
Environment=DISPLAY=:0.0 TTY=4
EnvironmentFile=/etc/X11/X.cfg
ExecStart=/usr/bin/Xorg vt${TTY} -s 0 -dpms -br -noreset -allowMouseOpenFail
ExecStartPost=/usr/lib/xorg/xorg-post-start
Restart=always

[Install]
WantedBy=graphical.target
