# packages/mpd/unit/mpd.service
[Unit]
Description=MPD 
ConditionPathExists=/etc/mpd.conf

[Service]
ExecStart=/usr/bin/mpd /etc/mpd.conf
Restart=always

[Install]
WantedBy=multi-user.target
