# packages/tvheadend/unit/tvheadend.service
[Unit]
Description=TV streaming server
After=syslog.target

[Service]
ExecStart=/usr/sbin/runtvheadend
Restart=always

[Install]
WantedBy=multi-user.target
