# packages/tvheadend/unit/tvheadend.service
[Unit]
Description=TV streaming server
After=syslog.target pvr.service
ConditionPathExists=/run/geexbox/pvr/tvheadend

[Service]
ExecStart=/usr/bin/tvheadend -c /etc/tvheadend -s
Restart=always

[Install]
WantedBy=multi-user.target
