# packages/tvheadend/unit/tvheadend.service
[Unit]
Description=TV streaming server
After=syslog.target

[Service]
EnvironmentFile=/etc/pvr.conf
ExecStartPre=/usr/bin/test ${BACKEND} = tvheadend
ExecStart=/usr/bin/tvheadend -c /etc/tvheadend -s
Restart=always

[Install]
WantedBy=multi-user.target
