# packages/net/openvpn/package/postinst
#!/bin/sh
echo "Post install:"
echo "  Remeber to enable client/server"
echo "  type the command:"
echo "  systemctl enable openvpn@<configuration>.service"
echo "  After please reboot"
if [ -f /tmp/cpopkg/client.conf ]; then
   cp -r /tmp/cpopkg/* /etc/openvpn/
   rm -rf /tmp/cpopkg
fi

echo "Done"

