# packages/busybox/unit/acpid.service
[Unit]
Description=ACPI Event Daemon
After=syslog.target

[Service]
ExecStart=/sbin/acpi-daemon
Restart=always

[Install]
WantedBy=multi-user.target
