# packages/busybox/unit/crond.service
[Unit]
Description=Command Scheduler
After=syslog.target

[Service]
ExecStart=/usr/sbin/crond -fS
Restart=always

[Install]
WantedBy=multi-user.target
