# packages/busybox/unit/klogd.service
[Unit]
Description=Kernel Logger Daemon
After=syslog.service

[Service]
ExecStart=/sbin/klogd -n
Restart=always

[Install]
WantedBy=syslog.target
