# packages/busybox/unit/klogd.service
[Unit]
Description=Kernel Logger Daemon
After=syslog.target

[Service]
ExecStart=/sbin/klogd -n
Restart=always

[Install]
WantedBy=multi-user.target
