# packages/busybox/unit/lspci.service
[Unit]
Description=Cache PCI devices

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/usr/bin/lspci > /tmp/pci'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
