# packages/busybox/unit/ntpd.service
[Unit]
Description=Network Time Protocol
After=network.target

[Service]
EnvironmentFile=/etc/ntpd.conf
ExecStart=/usr/sbin/ntpd $NTP_OPTS
Type=forking
Restart=always

[Install]
WantedBy=multi-user.target
