# packages/busybox/unit/ntpd.service
[Unit]
Description=Network Time Protocol
After=network.target

[Service]
ExecStartPre=/bin/sh -c 'while ! /bin/ping -c1 www.geexbox.org; do sleep 2; done'
EnvironmentFile=/etc/ntpd.conf
ExecStart=/usr/sbin/ntpd $NTP_OPTS
Type=forking
Restart=always

[Install]
WantedBy=multi-user.target
