# packages/busybox/unit/rdate.service
[Unit]
Description=rdate
After=connman.service

[Service]
Type=oneshot
ExecStart=/bin/rdate -s europe.pool.ntp.org
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
