# packages/busybox/unit/rdate.service
[Unit]
Description=rdate
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/runrdate
Restart=always
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
