# packages/busybox/unit/syslogd.service
[Unit]
Description=System Logger Daemon

[Service]
ExecStartPre=/bin/systemctl stop systemd-syslog-kmsg.service
ExecStart=/sbin/syslogd -n -C10240
Sockets=syslog.socket
Restart=always

[Install]
WantedBy=multi-user.target
