# packages/busybox/unit/syslogd.service
[Unit]
Description=System Logger Daemon
Before=klogd.service

[Service]
ExecStart=/sbin/syslogd -n -C10240
Restart=always

[Install]
WantedBy=syslog.target
