# packages/busybox/unit/telnetd@.service
[Unit]
Description=Telnet Per-Connection Server
After=syslog.target

[Service]
ExecStart=/usr/sbin/telnetd -i
StandardInput=socket
