# packages/busybox/unit/telnetd.socket
[Socket]
ListenStream=23
Accept=yes
EnvironmentFile=/etc/network
ExecStartPre=/usr/bin/test ${TELNET_SERVER} = true

[Install]
WantedBy=sockets.target
