# packages/multimedia/xbmc/profile.d/xbmc.sh
# XBMC settings

export XBMC_HOME="/usr/share/xbmc"
