# packages/libomxil-bellagio/init/omx.conf
description  "start OpenMAX Registrar"
author       "OpenBricks"

start on runlevel geexbox
task

script

  # start OpenMAX Bellagio plugins registrar
  omxregister-bellagio

  exit 0
end script
