# packages/udev/init/udevd.conf
description  "start udev daemon"
author       "OpenBricks"

start on startup 
respawn

exec udevd
