# packages/system/automountd/unit/automountd.service
[Unit]
Description=udev automounter daemon
Conflicts=installator.target
After=systemd-udev-settle.service

[Service]
ExecStart=/usr/bin/automountd.sh
Restart=always

[Install]
WantedBy=network.target
#WantedBy=multi-user.target
