# packages/networkmanager/package/postinst
#!/bin/sh

echo "Post install:"
echo "  Activating NetworkManager for next reboot ..."
systemctl enable NetworkManager
systemctl enable NetworkManager-dispatcher

