# packages/networkmanager/package/postinst
#!/bin/sh

echo "Post install:"
echo "  Activatin NetworkManager for next reboot ..."
systemctl enable NetworkManager

