# packages/snowball-init/unit/snowball-init.service
[Unit]
Description=Snowball init service
After=syslog.target

[Service]
ExecStart=/usr/local/sbin/snowball start
#ExecStop=/usr/local/sbin/snowball stop

[Install]
WantedBy=wpa_supplicant.service
