# packages/net/bluez/unit/bluetooth-setup.service
[Unit]
Description=Load Broadcom 43xx bluetooth firmware and enable HCI device
ConditionPathExists=/sys/class/mmc_host/mmc1

[Service]
Type=oneshot
ExecStart=/usr/bin/brcm43xx-bt-attach
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
