# packages/net/bluez/unit/bluetooth-setup.service
[Unit]
Description=Load bluetooth firmware and enable HCI device
ConditionPathExists=/sys/class/mmc_host/mmc1

[Service]
Type=oneshot
ExecStart=/usr/bin/bt-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
