# packages/showtime/unit/showtime.service
[Unit]
Description=HTS Showtime

[Service]
Environment=DISPLAY=:0.0
ExecStart=/usr/bin/showtime
Restart=always

[Install]
WantedBy=showtime.target
