# packages/net/polipo/unit/polipo.service
[Unit]
Description=Polipo Proxy Server

[Service]
ExecStart=/usr/bin/polipo
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=multi-user.target
