# packages/libs/hal/unit/haldaemon.service
[Unit]
Description=Hardware Manager
After=syslog.target acpid.service
Wants=acpid.service

[Service]
Type=dbus
BusName=org.freedesktop.Hal
ExecStart=/usr/sbin/hald --daemon=no --use-syslog

# We define not install section here, assuming that only legacy apps will
# require this, and hence on-demand activation should be fine.
