# packages/linux/init/bootdate.conf
description  "store boot date"
author       "GeeXboX"

start on started time
task

exec /bin/date > /tmp/bootdate
