# packages/linux/init/bootdate.conf
description	"store boot date"
author		"Davide Cavalca <davide@geexbox.org>"

start on started time
task

exec /bin/date > /tmp/bootdate
