# packages/xorg-server/init/xorg.conf
description  "start Xorg"
author       "OpenBricks"

start on mode hdtv and started persistent-data and runlevel geexbox and started video-card
stop on stopped enna
respawn

# tty used for the video display
env TTY=4

# the Xorg display
env DISPLAY=:0.0
export DISPLAY

pre-start script
  . /etc/funcs
  [ -f /var/data/config/xorg.conf ] && make_persistent file /etc/X11/xorg.conf /var/data/config/xorg.conf

  # Fix the /etc/X11 case-insensitive fs directory creation issue
  [ -d /etc/x11 ] && ln -s /etc/x11 /etc/X11

  exit 0
end script

exec Xorg vt$TTY -s 0 -dpms -br -noreset -allowMouseOpenFail

post-start script
  # wait a bit for Xorg to startup
  sleep 1

  # special hack for vmware to enhance screen resolution
  grep -q 'VMWARE(0)' /var/log/Xorg.0.log && VMWARE=yes

  . /etc/X11/X.cfg

  if [ $XORG_RESX = auto -a $XORG_RESY = auto -a $VMWARE = yes ]; then
    XORG_RESX=1280
    XORG_RESY=720
  fi

  [ $XORG_RESX != auto -a $XORG_RESY != auto ] && XRANDR_ARGS="-s ${XORG_RESX}x${XORG_RESY}"
  [ $XORG_RATE != auto ] && XRANDR_ARGS="$XRANDR_ARGS -r $XORG_RATE"

  [ -n "$XRANDR_ARGS" ] && xrandr -d $DISPLAY $XRANDR_ARGS >> /var/log/xrandr.log 2>&1 || true

  # set Xorg keymap
  MAP="us"
  grep qwerty /proc/cmdline && MAP="us"
  grep qwertz /proc/cmdline && MAP="de"
  grep azerty /proc/cmdline && MAP="fr"
  setxkbmap -display $DISPLAY $MAP
end script
