# packages/enna/init/enna-shutdown.conf
description  "shutdown the system at Enna exit"
author       "GeeXboX"

start on runlevel geexbox and stopped JOB=enna RESULT=ok
emits shutdown
task

exec initctl emit shutdown REBOOT=no
