# packages/misc/at/unit/atd.service
[Unit]
Description=Atd daemon

[Service]
ExecStart=/sbin/atd -d -f

[Install]
WantedBy=multi-user.target
