# packages/sound/alsa-utils/config/40-alsa.rules
# /lib/udev/rules.d/40-alsa.rules

# When a sound device is detected, restore the volume settings
KERNEL=="controlC[0-9]*", ACTION=="add", RUN+="/usr/lib/alsa-utils/alsa-mixer-setup"
