# packages/inadyn/unit/inadyn.service
[Unit]
Description=Inadyn a Simple Dynamic DNS client

[Service]
ExecStart=/usr/bin/inadyn

[Install]
WantedBy=multi-user.target
