# packages/multimedia/freevo/unit/freevo-recordserver.service
[Unit]
Description=Freevo recordserver daemon

[Service]
EnvironmentFile=/etc/freevo/recordserver.conf
ExecStart=/usr/bin/freevo $R_ON
ExecStop=/usr/bin/freevo $R_OFF
 
[Install]
WantedBy=freevo.target
