# packages/multimedia/freevo/unit/freevo.service
[Unit]
Description=Freevo

[Service]
Environment=DISPLAY=:0.0
ExecStart=/usr/bin/freevo
ExecStop=/usr/bin/freevo --stop
StandardOutput=syslog
Restart=always

[Install]
WantedBy=freevo.target
