# packages/net/redis/unit/redis-server.service
[Unit]
Description=redis-server
After=sysinit.target

[Service]
ExecStart=/usr/bin/redis-server
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
