# packages/multimedia/satip/unit/satip.service
[Unit]
Description=SatIP service

[Service]
ExecStart=/usr/bin/satip
Restart=always

[Install]
WantedBy=multi-user.target
