# packages/multimedia/dvbstreamer/unit/dvbstreamer.service
[Unit]
Description=dvbstreamer daemon
Before=freevo.service
ConditionPathExists=/root/.dvbstreamer

[Service]
Environment=HOME=/root
ExecStart=/usr/bin/dvbstreamer -a 0 -r -D -vv
Restart=always

[Install]
WantedBy=freevo.service
