# packages/MPlayer/unit/mplayer-video.service
[Unit]
Description=MPlayer video setup
After=lspci.service

[Service]
Type=oneshot
EnvironmentFile=/etc/video
ExecStart=/usr/lib/mplayer/setup-mplayer-video
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
