# packages/MPlayer/unit/mplayer-fonts.service
[Unit]
Description=MPlayer fonts setup

[Service]
Type=oneshot
ExecStart=/usr/lib/mplayer/setup-mplayer-fonts
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
