# packages/automounter-9p/unit/automounter-9p.service
[Unit]
Description=Mount 9p.2000u resources
ConditionPathExists=/etc/9p

[Service]
Type=oneshot
ExecStart=/usr/bin/automount_9p
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
