# packages/multimedia/mpd/unit/mpd.service
[Unit]
Description=MPD 
ConditionPathExists=/etc/mpd.conf

[Service]
ExecStart=/usr/bin/mpd --no-daemon /etc/mpd.conf
Restart=always

[Install]
WantedBy=multi-user.target
