# packages/x11/matchbox/unit/matchbox-window-manager.service
[Unit]
Description=Matchbox Window Manager

[Service]
ExecStart=/usr/bin/matchbox-window-manager -display :0.0 -use_titlebar no
Restart=always

[Install]
WantedBy=graphical.target
