# packages/net/bftpd/unit/bftpd.service
[Unit]
Description=FTP Server

[Service]
ExecStart=/usr/sbin/bftpd -D -c /etc/bftpd.conf

[Install]
WantedBy=multi-user.target
