# packages/brcm-patchram-plus/unit/brcm-patchram-plus.service
[Unit]
Description=Load bt firmware on cubox-i 

[Service]
Type=oneshot
ExecStart=/usr/bin/start-brcm-patchram-plus
RemainAfterExit=yes

[Install]
WantedBy=bluetooth.service
