# packages/brcm-patchram-plus/unit/brcm-patchram-plus.service
[Unit]
Description=Load bt firmware on cuboxi 

[Service]
Type=oneshot
ExecStart=/usr/bin/brcm_patchram_plus -d --patchram /lib/firmware/brcm/bcm4329.hcd --baudrate 3000000 --use_baudrate_for_download /dev/ttymxc3 --enable_hci

[Install]
WantedBy=bluetooth.service
