# packages/misc/boblight/unit/boblightd.service
[Unit]
Description=Boblightd service
After=graphical.target

[Service]
ExecStart=/usr/bin/boblightd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=xbmc.target
