# packages/misc/pvr/unit/pvr.service
[Unit]
Description=pvr backend

[Service]
Type=oneshot
EnvironmentFile=/etc/pvr.conf
ExecStartPre=/bin/mkdir -p /run/geexbox/pvr
ExecStart=/bin/touch /run/geexbox/pvr/${BACKEND}
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
