# packages/libs/libcec-daemon/unit/libcec-daemon.service
[Unit]
Description=Lib CEC Daemon

[Service]
ExecStart=/usr/bin/libcec-daemon
Restart=always

[Install]
WantedBy=multi-user.target
