# packages/multimedia/kodi/profile.d/kodi.sh
# XBMC settings

export KODI_HOME="/usr/share/kodi"
