# packages/automounter-smb/unit/automounter-smb.service
[Unit]
Description=Mount Samba/CIFS resources

[Service]
Type=oneshot
ExecStart=/usr/bin/automount_smb
ExecStop=/usr/bin/pre-umount
RemainAfterExit=yes

[Install]
WantedBy=network.target
