# packages/automounter-smb/unit/automounter-smb.service
[Unit]
Description=Mount Samba/CIFS resources
Requires=connman.service
After=connman.service

[Service]
Type=oneshot
ExecStart=/usr/bin/automount_smb
ExecStop=/usr/bin/pre-umount
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
