# packages/automounter-smb/unit/automounter-smb.service
[Unit]
Description=Mount Samba/CIFS resources

[Service]
Type=oneshot
ExecStart=/usr/bin/automount_smb
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
