# packages/net/btsync/unit/btsync-daemon.service
[Unit]
Description=btsync daemon
[Service]
Type=forking
ExecStartPre=/bin/sleep 20
ExecStart=/usr/bin/btsync --config /etc/btsync.conf
Restart=always
[Install]
WantedBy=multi-user.target
