# packages/vdr/init/vdr.conf
description  "configure and launch VDR"
author       "OpenBricks"

start on started modules and started persistent-data and runlevel geexbox
respawn

pre-start script
  . /etc/funcs

  make_persistent file /etc/wakeup /var/data/config/wakeup
  make_persistent file /etc/vdr/channels.conf /var/data/vdr/channels.conf create
  make_persistent file /etc/vdr/timers.conf /var/data/vdr/timers.conf create
  make_persistent file /etc/vdr/setup.conf /var/data/vdr/setup.conf create
  make_persistent dir /video /var/data/video create
  make_persistent dir /etc/vdr/plugins/epgsearch /var/data/vdr/epgsearch create
  make_persistent dir /etc/vdr/plugins/tvm2vdr /var/data/vdr/tvm2vdr create
end script

exec runvdr
