# packages/installator/unit/installator.service
[Unit]
Description=Installator

[Service]
Type=oneshot
ExecStart=/sbin/installator
StandardInput=tty
StandardOutput=inherit

[Install]
WantedBy=installator.target
