# packages/installator/unit/installator.target
[Unit]
Description=Installator
Requires=multi-user.target
Conflicts=rescue.service rescue.target
AllowIsolate=yes
