# packages/misc/mysql-server/unit/mysqld-root-pass-setup.service
[Unit]
Description=MySQL root password setup
After=mysqld.service
ConditionPathExists=!/etc/mysql/conf.d/root-pass.cnf

[Service]
Type=oneshot
ExecStart=/usr/lib/mysql/mysqld-root-pass-setup
RemainAfterExit=yes

[Install]
WantedBy=mysqld.service
