# packages/automounter-nfs/unit/automounter-nfs.service
[Unit]
Description=Mount NFS resources
ConditionPathExists=/etc/nfs
Requires=connman.service
After=connman.service

[Service]
Type=oneshot
ExecStart=/usr/bin/automount_nfs
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
