# packages/automounter-nfs/unit/automounter-nfs.service
[Unit]
Description=Mount NFS resources
ConditionPathExists=/etc/nfs

[Service]
Type=oneshot
ExecStart=/usr/bin/automount_nfs
RemainAfterExit=yes

[Install]
WantedBy=network.target
