# packages/pciutils/init/lspci.conf
description  "cache pci devices"
author       "GeeXboX"

start on started modules
task

exec lspci > /tmp/pci
