# packages/pciutils/init/lspci.conf
description  "cache pci devices"
author       "GeeXboX"

start on started modules
task

script
  lspci >/tmp/pci || exit 0
end script
