# packages/rtorrent/config/99-rtorrent.rules
KERNEL=="sd?1", ACTION=="add", RUN+="/bin/systemctl --no-block start  rtorrent-on-usb-disk@%k.service"