# packages/net/lighttpd/unit/lighttpd.service
[Unit]
Description=Lighttpd Daemon
ConditionPathExists=/run/geexbox/network/httpd
After=network.target

[Service]
ExecStartPre=/srv/www/htdocs/parse-www-config.sh
ExecStartPre=/usr/sbin/lighttpd -t -f /etc/lighttpd/lighttpd.conf
ExecStart=/usr/sbin/lighttpd -D -f /etc/lighttpd/lighttpd.conf

[Install]
WantedBy=multi-user.target