# packages/tvheadend/init/tvheadend.conf
description  "start Tvheadend"
author       "OpenBricks"

start on started modules
respawn

exec tvheadend -c /etc/tvheadend -C -s
