# packages/lirc/unit/lirc-setup.service
[Unit]
Description=LIRC initial setup
Before=lircd.service
ConditionPathExists=!/etc/lircd

[Service]
Type=oneshot
EnvironmentFile=/etc/remote
ExecStart=/usr/lib/lirc/setup-lirc

[Install]
WantedBy=lircd.service
