# packages/lirc/unit/lircd.service
[Unit]
Description=Lircd daemon
After=systemd-udevd.service
Before=xorg.service
ConditionPathExists=/etc/lirc/mycustom.conf

[Service]
Type=forking
EnvironmentFile=/etc/lirc/hardware.conf
ExecStartPre=/bin/mkdir -p /var/run/lirc/
ExecStart=/usr/sbin/lircd -H $DRIVER -d $DEVICE /etc/lirc/mycustom.conf -o /var/run/lirc/lircd

[Install]
WantedBy=multi-user.target
