# packages/lirc/unit/lircd.service
[Unit]
Description=LIRC daemon

[Service]
EnvironmentFile=/etc/lircd
ExecStart=/usr/bin/lircd -n -a --driver=$LIRC_DRIVER --device=$LIRC_DEVICE
Restart=on-failure

[Install]
WantedBy=multi-user.target
