# packages/multimedia/qtmediahub/unit/qtmediahub.service
[Unit]
Description=QtMediaHub

[Service]
Environment=DISPLAY=:0.0
ExecStart=/usr/share/qtmediahub/qtmediahub
Restart=always

[Install]
WantedBy=qtmediahub.target
