# packages/misc/cubox-led/unit/cubox-led.service
[Unit]
Description=Cubox led

[Service]
Type=oneshot
ExecStart=/usr/bin/cubox-led
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
