# packages/opkg/init/opkg.conf
description  "configure opkg"
author       "GeeXboX"

start on startup
task

script
  mkdir -p /var/lib/opkg

  exit 0
end script
