# packages/misc/serdisplib/unit/serdisplib.service
[Unit]
Description=SerdispLib LCD Daemon

[Service]
Type=forking
User=root
ExecStart=/usr/bin/l4m132c_tool -p /dev/usb/hiddev0 -t

[Install]
WantedBy=multi-user.target
