# packages/admin/installator/unit/installator.target
[Unit]
Description=Installator
Requires=basic.target
Conflicts=rescue.service rescue.target
AllowIsolate=yes
