# packages/net/brcm-patchram-plus/unit/brcm43xx-bt-init.service
[Unit]
Description=Load Broadcom 43xx bluetooth firmware
ConditionPathExists=/dev/ttymxc3

[Service]
Type=oneshot
ExecStart=/usr/bin/brcm43xx-bt-init
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
