# packages/net/brcm-patchram-plus/unit/brcm-patchram-plus.service
[Unit]
Description=Load Broadcom bluetooth firmware
ConditionPathExists=/dev/ttymxc3

[Service]
Type=oneshot
ExecStart=/usr/bin/start-brcm-patchram-plus
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
