# packages/network-backend/unit/network-link.service
[Unit]
Description=Network backend

[Service]
Type=oneshot
EnvironmentFile=/etc/network
ExecStartPre=/bin/mkdir -p /run/geexbox/network
ExecStart=/bin/touch /run/geexbox/network/${NETWORK_BACKEND}
RemainAfterExit=yes

[Install]
WantedBy=network.target
