# packages/udhcpc/unit/udhcpc.service
[Unit]
Description=Connection service (udhcpc)
After=syslog.target

[Service]
ExecStart=/usr/lib/udhcpc/net-udhcpc
RemainAfterExit=yes

[Install]
WantedBy=network.target
