# packages/freevo/unit/freevo.service
[Unit]
Description=Freevo

[Service]
Environment=DISPLAY=:0.0
ExecStart=/usr/bin/freevo
StandardOutput=syslog
Restart=always

[Install]
WantedBy=freevo.target
