# packages/net/samba-server/unit/nmbd.service
[Unit]
Description=nmbd Server
Before=smbd.service

[Service]
ExecStart=/usr/sbin/nmbd -F
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
