# packages/net/samba-server/unit/smbd.service
[Unit]
Description=Samba Server

[Service]
ExecStart=/usr/sbin/smbd -F
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
