# packages/net/samba-server/unit/smbd.service
[Unit]
Description=Samba File Sharing Service
ConditionPathExists=/run/geexbox/network/smbd
After=network.target

[Service]
ExecStart=/usr/sbin/smbd -F
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
