# packages/hdmiservice/unit/hdmi.service
[Unit]
Description=HDMI service
After=lspci.service

[Service]
ExecStart=/usr/bin/hdmistart
Restart=always

[Install]
WantedBy=graphical.target
