# packages/misc/reelbox-control/unit/ctrld.service
[Unit]
Description=reelbox-control

[Service]
ExecStart=/usr/sbin/reelbox-control
Restart=always

[Install]
WantedBy=multi-user.target
