# packages/pulseaudio/unit/pulsed.service
[Unit]
Description=Pulseaudio

[Service]
ExecStart=/usr/bin/pulseaudio --system --daemonize --high-priority --log-target=syslog --disallow-exit --disallow-module-loading=1

[Install]
WantedBy=multi-user.target
