# packages/misc/oscam/unit/oscam.service
[Unit]
Description=Oscam Services

[Service]
ExecStart=/usr/bin/oscam.sh
Restart=always

[Install]
WantedBy=multi-user.target
