# packages/multimedia/enna/unit/enna-setup.service
[Unit]
Description=Enna initial setup
Before=enna.service

[Service]
Type=oneshot
ExecStart=/usr/lib/enna/setup-enna
RemainAfterExit=yes

[Install]
WantedBy=enna.service
