# packages/multimedia/enna/unit/enna.target
[Unit]
Description=Enna Media Center
Requires=graphical.target
Conflicts=rescue.service rescue.target
After=graphical.target
AllowIsolate=yes
