# packages/multimedia/boxee/unit/boxee.service
[Unit]
Description=Boxee Media Center

[Service]
Environment=SDL_MOUSE_RELATIVE=0 DISPLAY=:0.0
ExecStart=/opt/boxee/Boxee
Restart=always

[Install]
WantedBy=boxee.target
