# packages/system/systemd/unit/hwclock-sync.service
[Unit]
Description=Synchronize System Clock With RTC

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/hwclock -s
ExecStop=/sbin/hwclock -w

[Install]
# Triggered by udev rule
