# packages/net/transmission/unit/transmission-daemon.service
[Unit]
Description=Transmission BitTorrent Daemon
After=network.target

[Service]
User=root
Type=notify
ExecStartPre=/sbin/sysctl -w net.core.rmem_max=8388608
ExecStartPre=/sbin/sysctl -w net.core.wmem_max=8388608
ExecStart=/usr/bin/transmission-daemon --log-error
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
