# packages/sgx540-OpenGLES/unit/pvr.service
[Unit]
Description=PowerVR SGX540 driver
Before=xorg.service
After=video-card.service
ConditionPathExists=/var/use_opengl

[Service]
Type=oneshot
ExecStartPre=/sbin/modprobe omap_gpu
ExecStart=/usr/bin/pvrsrvinit
RemainAfterExit=yes

[Install]
WantedBy=xorg.service
