# packages/automountd/unit/automountd.service
[Unit]
Description=udev automounter daemon
Conflicts=installator.target
After=udev-settle.service

[Service]
ExecStart=/usr/bin/automountd
Restart=always

[Install]
WantedBy=multi-user.target
