# packages/games/mame/unit/mame-config.service
[Unit]
Description=MAME configuration
ConditionPathExists=!/root/.mame/mame.ini

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'mkdir -p /root/.mame && ln -sf /etc/mame/ini/mame.ini /root/.mame/mame.ini'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.service
