# packages/dietsplash/unit/dietsplash-quit.service
[Unit]
Description=Terminate dietsplash Boot Screen
Before=xorg.service

[Service]
ExecStart=/usr/bin/splash.sh off
Type=oneshot
