# packages/dietsplash/unit/dietsplash.service
[Unit]
Description=Dietsplash Boot Screen
After=udev.service

[Service]
EnvironmentFile=/etc/dietsplash.conf
Type=oneshot
ExecStart=/usr/bin/splash.sh $ON
ExecStop=/usr/bin//splash.sh $OFF
RemainAfterExit=yes

[Install]
WantedBy=basic.target
