# packages/dietsplash/unit/dietsplash.service
[Unit]
Description=Dietsplash Boot Screen
After=udev.service

[Service]
Type=oneshot
ExecStart=/usr/bin/splash.sh on
RemainAfterExit=yes

[Install]
WantedBy=basic.target
