# packages/xcompmgr/unit/xcompmgr.service
[Unit]
Description=X composition manager
After=xorg.service

[Service]
Environment=DISPLAY=:0.0
ExecStart=/usr/bin/xcompmgr -n
Restart=always

[Install]
WantedBy=graphical.target
