# packages/net/openssh/unit/sshd.service
[Unit]
Description=SSH Per-Connection Server
After=syslog.target

[Service]
ExecStart=/usr/bin/opensshkeygen

[Install]
WantedBy=network.target
