# packages/system/systemd/config/wake-on-lan.service
[Unit]
Description=Wake On Lan Service

[Service]
Type=oneshot
EnvironmentFile=/etc/network
ExecStart=/usr/sbin/wakeonlan.sh
RemainAfterExit=yes

[Install]
WantedBy=network.target

