# packages/system/busybox/unit/acpid.service
[Unit]
Description=ACPI Event Daemon
After=syslog.target
ConditionPathExists=/sbin/acpid

[Service]
ExecStart=/sbin/acpid -d
Restart=always

[Install]
WantedBy=multi-user.target
