# packages/system/busybox/unit/logrotate.service
[Unit]
Description=Log Rotate for optimize space

[Service]
ExecStart=/bin/logrotate.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
