# packages/system/busybox/unit/telnetd.socket
[Unit]
Description=Telnet socket
ConditionPathExists=/run/geexbox/network/telnetd
After=network-link.service

[Socket]
ListenStream=23
Accept=yes

[Install]
WantedBy=sockets.target
