# packages/rtorrent/unit/rtorrent.service
[Unit]
Description=rtorrent daemon

[Service]
ExecStart=/usr/bin/runrtorrent
#User=nobody
#Group=nobody
Nice=10
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
