# packages/mysql-server/unit/mysqld-bootstrap.service
[Unit]
Description=MySQL database bootstrap
Before=mysqld.service
ConditionPathExists=!/var/lib/mysql/mysql

[Service]
Type=oneshot
ExecStart=/usr/lib/mysql/mysqld-bootstrap
RemainAfterExit=yes

[Install]
WantedBy=mysqld.service
