# packages/systemd/unit/hwclock-sync.service
[Unit]
Description=Synchronize System Clock With RTC
After=sysinit.target

[Service]
Type=forking
RemainAfterExit=yes
TimeoutSec=0
ExecStart=/sbin/hwclock -s
ExecStop=/sbin/hwclock -w
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
