# packages/systemd/unit/tmp.mount
[Unit]
Description=Temporary Directory
Before=local-fs.target

[Mount]
What=tmpfs
Where=/tmp
Type=tmpfs
Options=mode=1777
