# packages/tegra-libs/unit/nvrm.service
[Unit]
Description=nVidia Tegra daemon
Before=xorg.service
After=video-card.service
ConditionPathExists=/var/use_opengl

[Service]
ExecStartPre=/bin/chmod 0666 /sys/power/state
ExecStart=/usr/sbin/nvrm_daemon
Restart=always

[Install]
WantedBy=xorg.service
