# packages/tegra-libs/unit/tegra-setup.service
[Unit]
Description=nVidia Tegra initalization
Before=xorg.service
After=video-card.service

[Service]
Type=oneshot
ExecStart=/usr/lib/tegra-libs/tegra-setup

[Install]
WantedBy=xorg.service
