# packages/drivers/eventlircd/unit/irexec.service
[Unit]
Description=LIRC irexec Daemon
Requires=eventlircd.service
After=eventlircd.service

[Service]
Type=forking
EnvironmentFile=/etc/lirc/irexec
User=root
ExecStart=/usr/bin/irexec --daemon $IREXEC_OPTS

[Install]
WantedBy=multi-user.target
