# packages/misc/weston/unit/weston.service
[Unit]
Description=Weston compositor
After=lspci.service

[Service]
Environment=XDG_RUNTIME_DIR=/tmp/
ExecStart=/usr/bin/weston --backend=drm-backend.so
Restart=always

[Install]
WantedBy=graphical.target
