# packages/admin/plymouth/unit/plymouth-quit-wait.service
[Unit]
Description=Wait for Plymouth Boot Screen to Quit
After=rc-local.service plymouth-start.service systemd-user-sessions.service

[Service]
ExecStart=-/bin/plymouth --wait
Type=oneshot
TimeoutSec=20
