# packages/admin/plymouth/unit/plymouth-kexec.service
[Unit]
Description=Show Plymouth Reboot with kexec Screen
After=getty@tty1.service prefdm.service plymouth-start.service
Before=kexec.service
DefaultDependencies=no
ConditionKernelCommandLine=!plymouth.enable=0

[Service]
ExecStart=/sbin/plymouthd --mode=shutdown --attach-to-session
ExecStartPost=-/bin/plymouth show-splash
Type=forking
