# packages/admin/plymouth/unit/systemd-ask-password-plymouth.service
[Unit]
Description=Forward Password Requests to Plymouth
Documentation=http://www.freedesktop.org/wiki/Software/systemd/PasswordAgents
DefaultDependencies=no
Conflicts=shutdown.target
After=plymouth-start.service
Before=shutdown.target
ConditionKernelCommandLine=!plymouth.enable=0
ConditionPathExists=/run/plymouth/pid

[Service]
ExecStart=/bin/systemd-tty-ask-password-agent --watch --plymouth
