# packages/admin/plymouth/unit/systemd-ask-password-plymouth.path
[Unit]
Description=Forward Password Requests to Plymouth Directory Watch
Documentation=http://www.freedesktop.org/wiki/Software/systemd/PasswordAgents
DefaultDependencies=no
Conflicts=shutdown.target
After=plymouth-start.service
Before=basic.target shutdown.target
ConditionKernelCommandLine=!plymouth.enable=0
ConditionPathExists=/run/plymouth/pid

[Path]
DirectoryNotEmpty=/run/systemd/ask-password
MakeDirectory=yes
