# packages/multimedia/transmission/unit/transmission-daemon.service
[Unit]
Description=Transmission BitTorrent Daemon
After=network.target

[Service]
User=root
Type=notify
ExecStart=/usr/bin/transmission-daemon --log-error
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
