# packages/e2fsprogs/unit/sdcard-resize.service
[Unit]
Description=Resize and enlarge sdcard

[Service]
Type=oneshot
ExecStart=/usr/sbin/resize-and-enlarge

[Install]
WantedBy=multi-user.target
