# packages/sgx530-OpenGLES/unit/pvr.service
[Unit]
Description=PowerVR SGX530 driver
Before=xorg.service
After=video-card.service
ConditionPathExists=/var/use_opengl

[Service]
Type=oneshot
ExecStart=/usr/bin/do-pvr-init
RemainAfterExit=yes

[Install]
WantedBy=xorg.service
