# packages/misc/cups/unit/cups.service
[Unit]
Description=The CUPS scheduler
After=syslog.target

[Service]
Type=forking
ExecStartPre=/usr/bin/cups-fix-perms
ExecStartPre=/usr/sbin/cupsd -t
ExecStart=/usr/sbin/cupsd
ExecStartPost=/usr/bin/udevadm trigger --subsystem-match=usb --attr-match=bInterfaceClass=07 --attr-match=bInterfaceSubClass=01 --action=add
ExecStartPost=/usr/bin/udevadm trigger --subsystem-match=usb --property-match=DEVNAME="/dev/usb/lp*" --action=add

[Install]
WantedBy=multi-user.target
