# packages/multimedia/vdr-lcdproc/unit/lcdproc.service
[Unit]
Description=LCD Proc
#Requires=xorg.service vdr.service
#After=lircd.service

[Service]
ExecStart=/usr/sbin/LCDd -f -c /etc/LCDd.conf
Restart=always

[Install]
WantedBy=multi-user.target
