# packages/net/dropbear/unit/sshd@.service
[Unit]
Description=SSH Per-Connection Server
After=syslog.target

[Service]
EnvironmentFile=/etc/dropbear/dropbear.options
ExecStartPre=/usr/bin/dropbearkeygen
ExecStart=/usr/bin/dropbear -i $OPTIONS
StandardInput=socket
KillMode=process