# packages/net/dropbear/unit/sshd.socket
[Unit]
Description=SSH socket
ConditionPathExists=/run/geexbox/network/sshd
After=network-link.service

[Socket]
ListenStream=22
Accept=yes

[Install]
WantedBy=sockets.target
