# packages/python/pyload/unit/pyload.service
[Unit]
Description=pyload service

[Service]
Type=oneshot
ExecStart=/usr/bin/pyload
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
