# packages/ti-dsp/init/dsp.conf
description  "load DSP modules"
author       "OpenBricks"

start on started modules
task

script
  modprobe bridgedriver base_img=/lib/dsp/baseimage.dof
  modprobe dspbridge

  exit 0
end script
