# packages/multimedia/kodi/unit/silent.service
[Unit]
Description=Silent messages
After=kodi.service

[Service]
Type=oneshot
ExecStart=/usr/bin/silent.sh

[Install]
WantedBy=kodi.service
