# packages/xbmc/unit/silent.service
[Unit]
Description=Silent messages
Before=xbmc.service

[Service]
Type=oneshot
ExecStart=/usr/bin/silent.sh

[Install]
WantedBy=xbmc.service
