# packages/xbmc/unit/xbmc-pvr.target
[Unit]
Description=XBMC Media Center pvr edition
Requires=graphical.target
Conflicts=rescue.service rescue.target
After=graphical.target
AllowIsolate=yes
