# packages/xbmc/unit/xbmc.service
[Unit]
Description=XBMC Media Center

[Service]
Environment=DISPLAY=:0.0 HOME=/root
ExecStart=/usr/lib/xbmc/xbmc.bin --standalone -fs
Restart=always

[Install]
WantedBy=xbmc.target
