# packages/xbmc/unit/xbmc-setup.service
[Unit]
Description=XBMC Early Setup
Before=xbmc.service

[Service]
Type=oneshot
ExecStart=/usr/lib/xbmc/xbmc-default-settings.sh

[Install]
WantedBy=xbmc.service
