# packages/xbmc/unit/xbmc.service
[Unit]
Description=XBMC Media Center
Requires=xorg.service
After=xorg.service

[Service]
Environment=DISPLAY=:0.0 SDL_MOUSE_RELATIVE=0 HOME=/root
EnvironmentFile=/etc/profile
ExecStart=/usr/bin/runxbmc --standalone -fs -l /var/run/lirc/lircd
Restart=on-failure
# KillSignal=SIGKILL

[Install]
WantedBy=xbmc.target
