# packages/multimedia/kodi-wiiremote/unit/kodi-wiiremote.service
[Unit]
Description=wiiremote for kodi
After=bluetooth.service kodi.service

[Service]
Type=oneshot
ExecStart=/usr/bin/start-wiiremote
RemainAfterExit=yes

[Install]
WantedBy=bluetooth.target
