# packages/misc/cubox-hdmi/unit/hdmi.service
[Unit]
Description=HDMI service
After=lspci.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c /usr/bin/hdmiwait


[Install]
WantedBy=xorg.service
