# packages/alsa-utils/unit/mixer.service
[Unit]
Description=ALSA mixer configuration
After=sysinit.target

[Service]
Type=oneshot
EnvironmentFile=/etc/audio
ExecStart=/usr/lib/alsa-utils/alsa-mixer-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
