# packages/xbmc-eventclient/unit/wii4xbmc.service
[Unit]
Description=wiiremote for xbmc
After=xbmc.service
Requires=bluetooth.service

[Service]
Type=oneshot
ExecStart=/usr/bin/start-wiiremote
RemainAfterExit=yes

[Install]
WantedBy=xbmc.service
