# packages/net/mysql-server/unit/mysqld.service
[Unit]
Description=MySQL database server

[Service]
ExecStart=/usr/bin/mysqld_safe
ExecStop=/usr/bin/mysqladmin shutdown
ExecReload=/usr/bin/mysqladmin reload

[Install]
WantedBy=multi-user.target
