# packages/busybox/profile.d/system.sh
# System Settings

export PATH="/bin:/sbin:/usr/bin:/usr/sbin"
export HOSTNAME=$(cat /etc/hostname)

alias ll='ls -l'
