# packages/admin/dietsplash/unit/dietsplash-quit.service
[Unit]
Description=Terminate dietsplash Boot Screen
Before=xorg.service

[Service]
EnvironmentFile=/etc/dietsplash.conf
ExecStart=/usr/bin/splash.sh $OFF
Type=oneshot
