# packages/automounter-9p/init/9p.conf
description  "mount 9p.2000u resource sharing"
author       "OpenBricks"

start on started networking
task

pre-start exec test -f /etc/9p -a -f /etc/network

exec automount_9p
