# packages/drivers/ralink-8070-3070-3370-5372-usb/config/70-ralinkusb.rules
# /lib/udev/rules.d/70-ralinkusb.rules

ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="2001", ATTR{idProduct}=="3c15", RUN+="/sbin/modprobe -qba rt5370sta"
