# packages/syslink/unit/syslink.service
[Unit]
Description=Ducati Syslink daemon
ConditionPathExists=/etc/syslink.conf
After=syslog.target

[Service]
EnvironmentFile=/etc/syslink.conf
Type=forking
ExecStart=/usr/bin/syslink_daemon.out $SYSM3_IMAGE $APPM3_IMAGE
Restart=always

[Install]
WantedBy=multi-user.target
