# packages/multimedia/MPlayer/unit/mplayer-audio.service
[Unit]
Description=MPlayer audio setup
After=lspci.service

[Service]
Type=oneshot
EnvironmentFile=/etc/audio
ExecStart=/usr/lib/mplayer/setup-mplayer-audio
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
