# packages/misc/irqbalance/unit/irqbalance.service
[Unit]
Description=irqbalance
After=sysinit.target

[Service]
ExecStart=/usr/sbin/irqbalance
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
