# packages/boxee/unit/boxee.service
[Unit]
Description=Boxee Media Center

[Service]
Environment=DISPLAY=:0.0
ExecStart=/opt/boxee/Boxee
Restart=always

[Install]
WantedBy=boxee.target
