# packages/boxee/unit/boxee.service
[Unit]
Description=XBMC Media Center

[Service]
Environment=DISPLAY=:0.0
ExecStartPre=/usr/bin/python /usr/share/xbmc/FEH.py
ExecStart=/usr/lib/xbmc/xbmc.bin
Restart=always

[Install]
WantedBy=xbmc.target
