# packages/dvbstreamer/unit/dvbstreamer.service
[Unit]
Description=dvbstreamer daemon
Before=freevo.service
ConditionPathExists=/root/.dvbstreamer

[Service]
Type=forking
Environment=HOME=/root
ExecStart=/usr/bin/dvbstreamer -a 0 -r -d -D

[Install]
WantedBy=freevo.service
