# packages/networkmanager/unit/nm.service
[Unit]
Description=Create basic config for NetworkManager service
After=syslog.target network-link.service
ConditionPathExists=/run/geexbox/network/networkmanager

[Service]
ExecStart=/usr/bin/nm-parse-configuration

[Install]
WantedBy=network.target
