# packages/net/shadowsocks/unit/ss-local.service
[Unit]
Description=ShadowSocks Proxy Client

[Service]
Type=simple
PIDFile=/var/run/ss-local.pid
ExecStart=/usr/bin/ss-local -c /etc/shadowsocks/ss-local.json -f /var/run/ss-local.pid

[Install]
WantedBy=multi-user.target
