# packages/net/shadowsocks/unit/ss-server.service
[Unit]
Description=ShadowSocks Proxy Server

[Service]
Type=simple
PIDFile=/var/run/ss-server.pid
ExecStart=/usr/bin/ss-server -c /etc/shadowsocks/ss-server.json -f /var/run/ss-server.pid

[Install]
WantedBy=multi-user.target
