# packages/connman/unit/connman.service
[Unit]
Description=Connection service
After=syslog.target network-link.service
ConditionPathExists=/run/geexbox/network/connman

[Service]
Type=dbus
BusName=net.connman
ExecStartPre=/usr/lib/connman/connman-parse-configuration
ExecStart=/usr/sbin/connmand -n

[Install]
WantedBy=network.target
