# packages/connman/unit/mount-a.service
Unit]
Description=Mount all devices in /etc/fstab
After=connman.service
Before=graphical.target

[Service]
Type=oneshot
ExecStart=/bin/mount -a
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
